# infra/ansible/templates/gossip-logrotate.service.j2
[Unit]
Description=Rotate gossip observer logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/gossip
